# services/chatbot/run
#!/bin/bash
export USER=app

su - node -c 'cd /src; npm start'
